# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=ca
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
silent=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:grepmarx]
command=gunicorn -w 3 --bind 0.0.0.0:8001 run:app 
#command=gunicorn -w 3 --bind unix:grepmarx.sock run:app 
autostart=true

[program:celery_worker]
command=celery -A celery_worker.celery worker --loglevel=info
autostart=true